# systemd/foreign-docker-systemd.service.j2
[Unit]
Description={{ vendor.name }} {{ box }} {{ stage }}

[Service]
{% if host.docker is defined and host.docker.url is defined %}
Environment=DOCKER_HOST={{ host.docker.url }}
{% endif %}
User=root
Group=root
ExecStart=/opt/{{ system.service.name }}/bin/start
ExecStop=/opt/{{ system.service.name }}/bin/stop
Restart=on-failure
RestartSec=20
TimeoutStopSec=180

[Install]
WantedBy=multi-user.target
