# systemd/rootless-systemd.service.j2
[Unit]
Description={{ vendor.name }} {{ box }} {{ stage }}
Requires=docker.service
After=docker.service

[Service]
Environment={{ vendor.short.name|upper }}_{{ box|upper }}_BASE_DIR={{ host.base.dir }}
{% if host.docker is defined and host.docker.url is defined %}
Environment=DOCKER_HOST={{ host.docker.url }}
{% endif %}
ExecStart={{ host.base.dir }}/bin/start
ExecStop={{ host.base.dir }}/bin/stop
Restart=on-failure
RestartSec=20
TimeoutStopSec=180

[Install]
WantedBy=default.target
