# systemd/rootful-hw-bw.service
[Unit]
Description=hexway bw prod
Requires=docker.service
After=docker.service

[Service]
User=hw-bw-srv
Group=hw-bw-srv
ExecStart=/opt/hw-bw/bin/start
ExecStop=/opt/hw-bw/bin/stop
Restart=on-failure
RestartSec=20
TimeoutStopSec=180

[Install]
WantedBy=multi-user.target